# pipeline/cloud-run/com.realzaragoza.sofascore-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.realzaragoza.sofascore-weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/robertsoliva/Desktop/Projects/Real-Zaragoza/pipeline/cloud-run/run_weekly_sofascore.sh</string>
  </array>

  <!-- Every Tuesday at 07:30 local time -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>2</integer>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/sofascore_launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sofascore_launchd.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
